# Pinak_Services/memory_service/scripts/com.pinak.memory.backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.pinak.memory.backup</string>

    <key>ProgramArguments</key>
    <array>
      <string>/bin/bash</string>
      <string>/Users/abhi-macmini/clawd-simba/Pinak_Projects/Pinak_Services/memory_service/scripts/pinak-memory-backup.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
      <key>Hour</key>
      <integer>3</integer>
      <key>Minute</key>
      <integer>15</integer>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/abhi-macmini/Library/Logs/pinak-memory-backup.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/abhi-macmini/Library/Logs/pinak-memory-backup.err.log</string>
  </dict>
</plist>
